# internet_led
#!/bin/sh /etc/rc.common

# نام سرویس شما
USE_PROCD=1
START=99
STOP=15

# مسیر اسکریپت اصلی که باید اجرا شود
SERVICE_SCRIPT="/usr/bin/internet_led_status.sh"

# تابع start_service: توسط procd برای شروع سرویس فراخوانی می‌شود
start_service() {
    # procd یک نمونه جدید از سرویس را ایجاد می‌کند
    procd_open_instance

    # procd دستور اصلی را برای اجرا تنظیم می‌کند
    # اسکریپت شما با /bin/sh اجرا خواهد شد
    procd_set_param command /bin/sh "$SERVICE_SCRIPT"

    # این بخش کلیدی است:
    # اگر اسکریپت به هر دلیلی (مثلاً خطا) بسته شود، procd به صورت خودکار آن را مجدداً راه‌اندازی می‌کند.
    # پارامتر اول (threshold): تعداد تلاش برای راه‌اندازی مجدد در یک بازه زمانی.
    # پارامتر دوم (timeout): اگر در این بازه زمانی (مثلاً 60 ثانیه) بیش از 5 بار تلاش ناموفق بود، متوقف شود.
    # پارامتر سوم (delay): فاصله زمانی (به ثانیه) بین هر تلاش برای راه‌اندازی مجدد.
    procd_set_param respawn ${respawn_threshold:-3600} ${respawn_timeout:-5} ${respawn_delay:-5}

    # بستن نمونه procd و اجرای نهایی
    procd_close_instance
}

# تابع stop_service (اختیاری اما توصیه شده):
# procd به صورت خودکار پروسه را متوقف می‌کند، اما این تابع برای تعریف صریح وجود دارد.
# stop_service() {
#     # این بخش به procd می‌گوید که چگونه سرویس را متوقف کند (معمولاً با ارسال سیگنال TERM).
#     # در اکثر موارد، این تابع نیاز به کد خاصی ندارد و procd کار را انجام می‌دهد.
# }
